# debian/postinst
#!/bin/sh
# postinst script for #PACKAGE#
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)
        # TODO actually install these on postinst? or download them as part
        # of building the package and include the files, installing them with
        # the rest of our python scripts in /usr/share/amppki/.
        echo "=============================================================="
        echo "Debian has really old python packages, if you haven't already:"
        echo
        echo "pip install pyramid pyramid-assetviews pyramid-chameleon \\"
        echo "\tpyramid-mako pyasn1 pyasn1-modules waitress"
        echo "=============================================================="

        # the first thing the user will need to do is initialise the CA
        # directory and we don't let them configure it, so just do it now
        if [ ! -d /etc/amppki/ ]; then
            /usr/sbin/ampca initialise
        fi

        if [ ! -f /usr/share/amppki/production.ini ]; then
            cp /usr/share/amppki/production.ini-dist /usr/share/amppki/production.ini
        fi

        # enable the SSL module
        /usr/sbin/a2enmod ssl

        if apache2ctl configtest 2>/dev/null; then
            #invoke-rc.d apache2 force-reload || true
            # use "restart" as a workaround for bug #558608
            invoke-rc.d apache2 restart || true
        else
            echo "Your apache2 configuration is broken, please fix it and restart apache2 manually."
        fi
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
