# debian/prerm
#!/bin/sh
# prerm script for amplet2
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <prerm> `remove'
#        * <old-prerm> `upgrade' <new-version>
#        * <new-prerm> `failed-upgrade' <old-version>
#        * <conflictor's-prerm> `remove' `in-favour' <package> <new-version>
#        * <deconfigured's-prerm> `deconfigure' `in-favour'
#          <package-being-installed> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    remove)
        # TODO use dh_apache2
        RESTART=false
        if [ -d /etc/apache2/conf-available/ ]; then
            # try to uninstall the config for apache2 in jessie or similar
            if [ -e /etc/apache2/conf-available/amppki.conf ]; then
                /usr/sbin/a2disconf -q amppki
                rm /etc/apache2/conf-available/amppki.conf
                RESTART=true
            fi
        elif [ -d /etc/apache2/conf.d/ ]; then
            # try to uninstall the config for apache2 in wheezy or similar
            if [ ! -e /etc/apache2/conf.d/amppki.conf ]; then
                rm /etc/apache2/conf.d/amppki.conf
                RESTART=true
            fi
        else
            echo "ERROR: Couldn't determine apache2 configuration directory"
            exit 1
        fi

        # restart if we disabled the config
        if $RESTART; then
            if apache2ctl configtest 2>/dev/null; then
                #invoke-rc.d apache2 force-reload || true
                # use "restart" as a workaround for bug #558608
                invoke-rc.d apache2 restart || true
            else
                echo "Your apache2 configuration is broken, please fix it and restart apache2 manually."
            fi
        fi
    ;;

    deconfigure|upgrade)
    ;;

    failed-upgrade)
    ;;

    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
